# docker/spark/supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/opt/spark/sbin/start-thriftserver.sh --master %(ENV_SPARK_MASTER)s
